# systemd/elog_db_update.service
[Unit]
Description=Will generate graphs from RRD file

[Service]
Type=simple
ExecStart=/usr/bin/php /usr/local/sbin/elog_db_update.php
