# systemd/gen_graphs_month.service
[Unit]
Description=Will generate graphs from RRD file

[Service]
Type=simple
ExecStart=/usr/bin/nice -n 20 /usr/local/sbin/gen_graphs.sh month
