# systemd/generate_graphs_1h.service
[Unit]
Description=Will generate graphs from RRD file

[Service]
Type=simple
ExecStart=/usr/local/sbin/graphs_create.sh 1
